# systemd/mpclipboard-client.service
[Unit]
Description=MPClipboard Client

[Service]
ExecStart=mpclipboard-client
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s
